# .config/systemd/user/feh-bing.timer
[Unit]
Description=Run feh-bing with delay

[Timer]
OnBootSec=1min
Unit=feh-bing.service

[Install]
WantedBy=timers.target